# .docker/deployments/gunicorn.conf
[program:gunicorn]
command=/home/user/app/scripts/run_gunicorn.sh
directory=/home/user/app
user=user
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile =/logs/gunicorn.log
redirect_stderr=true
